# saltstack/salt/backup/hadduck@.timer
[Unit]
Description=Daily hadduck remote backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
