# saltstack/salt/backup/hadduck@.timer
[Unit]
Description=Daily hadduck remote backup

[Timer]
OnCalendar=*-*-* 00:10:00

[Install]
WantedBy=timers.target
